# deploy/systemd/sonar-daily-curves.service
[Unit]
Description=SONAR Daily Curves (NSS sovereign yields)
Documentation=file:///home/macro/projects/sonar-engine/docs/ops/systemd-deployment.md
After=network-online.target
Wants=network-online.target
StartLimitBurst=2
StartLimitIntervalSec=3600

[Service]
Type=oneshot
User=macro
Group=macro
WorkingDirectory=/home/macro/projects/sonar-engine
EnvironmentFile=/home/macro/projects/sonar-engine/.env
# CAL-138 ship — --all-t1 iterates US + 6 EA members (DE/PT/IT/ES/FR/NL).
# Countries without a wired curve path (PT/IT/ES/FR/NL per
# CAL-CURVES-EA-PERIPHERY) skip with a warning; the run exits 0 if at
# least one country persists. Anchor on yesterday so EOD FRED / ECB SDW /
# Bundesbank / TE data is settled.
ExecStart=/bin/bash -lc '/home/macro/.local/bin/uv run python -m sonar.pipelines.daily_curves --all-t1 --date $(date -u -I -d yesterday)'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
Restart=on-failure
RestartSec=600

[Install]
WantedBy=multi-user.target
