# deploy/systemd/sonar-daily-curves.timer
[Unit]
Description=SONAR Daily Curves Timer (06:00 UTC)
Documentation=file:///home/macro/projects/sonar-engine/docs/ops/systemd-deployment.md

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
Unit=sonar-daily-curves.service

[Install]
WantedBy=timers.target
