# deploy/systemd/sonar-daily-monetary-indices.service
[Unit]
Description=SONAR Daily Monetary Indices (M1/M2/M3/M4)
Documentation=file:///home/macro/projects/sonar-engine/docs/ops/systemd-deployment.md
After=network-online.target sonar-daily-curves.service sonar-daily-overlays.service
Wants=network-online.target
# M3 reads EXPINF (overlays); M1/M2 read NSS forwards (curves).
StartLimitBurst=2
StartLimitIntervalSec=3600

[Service]
Type=oneshot
User=macro
Group=macro
WorkingDirectory=/home/macro/projects/sonar-engine
EnvironmentFile=/home/macro/projects/sonar-engine/.env
ExecStart=/bin/bash -lc '/home/macro/.local/bin/uv run python -m sonar.pipelines.daily_monetary_indices --all-t1 --date $(date -u -I -d yesterday)'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
Restart=on-failure
RestartSec=600

[Install]
WantedBy=multi-user.target
