# deploy/systemd/sonar-daily-financial-indices.timer
[Unit]
Description=SONAR Daily Financial Indices Timer (07:00 UTC)
Documentation=file:///home/macro/projects/sonar-engine/docs/ops/systemd-deployment.md

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
Unit=sonar-daily-financial-indices.service

[Install]
WantedBy=timers.target
